# etc/chantek.service
[Unit]
Description=Chantek

[Service]
Type=simple
ExecStart=/usr/bin/python /var/lib/chantek/server.py

[Install]
WantedBy=multi-user.target