# etc/ubuntu-init.conf
description "Chantek"
start on runlevel [2345]
stop on runlevel [!2345]
respawn limit 3 10
chdir /var/lib/chantek
exec python server.py
